# deploy/amenity-booking.service
[Unit]
Description=Amenity Booking Scheduler
After=network-online.target
Wants=network-online.target
# Disable the start-rate limiter so Restart=always is never throttled.
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/amenity-automation
# xvfb-run gives Chrome a virtual display so the existing headless:false launch works unchanged.
ExecStart=/usr/bin/xvfb-run -a /usr/bin/node reserve.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
